# etc/update-motd.d/10-sysinfo
#!/bin/bash
echo "Host: $(hostname -f)"
echo "OS: $(lsb_release -d | cut -f2)"
echo "Kernel: $(uname -r)"
