# etc/update-motd.d/20-uptime
#!/bin/bash
echo "Uptime: $(uptime -p) | Load: $(uptime | awk -F'load average:' '{print $2}')"
echo
